# workspaces-default-settings-1.0/etc/xdg/autostart/workspaces-app-folders.desktop
[Desktop Entry]
Type=Application
Encoding=UTF-8
Name=Workspaces App Folders
Comment=Create Workspaces application folders
Exec=workspaces-app-folders
Terminal=false
Icon=folder
NoDisplay=true
StartupNotify=false
Categories=System;Settings;
